# supvisors/test/use_cases/scenario_2/template_etc/console/group_hci.ini
[group:scen2_hci]
programs=chart_view,sensor_control,sensor_view,hci_internal_check_data_bus,hci_internal_data_bus

[program:chart_view]
command=./bin/config_manager.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:sensor_control]
command=./bin/sensor_control.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:sensor_view]
command=./bin/sensor_view.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:hci_internal_check_data_bus]
command=./bin/internal_check_data_bus.sh %(program_name)s
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:hci_internal_data_bus]
command=./bin/internal_data_bus.sh %(program_name)s
autostart=false
startsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

