# supvisors/test/use_cases/scenario_3/template_etc/console/programs_console.ini
[program:scen3_chart_view]
command=./bin/scen3/chart_view.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_item_control]
command=./bin/scen3/item_control.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_check_internal_data_bus]
command=./bin/scen3/check_internal_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5
