# supvisors/test/etc/my_movies/program_web_server.ini
[program:web_server]
command=python ./scripts/process_app.py -n %(program_name)s
autostart=false
startsecs=5
startretries=0
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
