# supervisors/test/etc/test/program_xfontsel.ini
[program:xfontsel]
command=/usr/bin/xfontsel
process_name=xfontsel
autostart=false
autorestart=false
startsecs=1
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
