# supervisors/test/etc/test/program_sleep.ini
[program:sleep]
command=/bin/sleep 5
process_name=sleep
autostart=false
autorestart=false
startsecs=1
exitcodes=0
stopwaitsecs=2
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
