# supervisors/test/etc/test/program_segv.ini
[program:segv]
command=sleep 1 ; exit 11
process_name=segv
autostart=false
autorestart=unexpected
startsecs=2
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
