# supervisors/test/etc/my_movies/program_disk_reader.ini
[program:disk_reader]
command=python ./scripts/process_app.py -n %(program_name)s -x 20
process_name=disk_%(process_num)02d_reader
numprocs=3
numprocs_start=1
autostart=false
autorestart=false
startsecs=2
startretries=0
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
