# supervisors/test/etc/my_movies/program_error_disk_writer.ini
[program:error_disk_writer]
command=python ./scripts/process_app.py -n %(program_name)s -x 2
autostart=false
autorestart=false
startsecs=5
startretries=1
exitcodes=1
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
