# supvisors/test/use_cases/scenario_3/etc/server/group_server.ini
[group:scen3_srv]
programs=scen3_item_manager,scen3_track_manager,scen3_system_health,scen3_common_bus_interface,scen3_check_common_data_bus,scen3_check_internal_data_bus

[program:scen3_item_manager]
command=./bin/scen3/item_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_track_manager]
command=./bin/scen3/track_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_system_health]
command=./bin/scen3/system_health.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_check_common_data_bus]
command=./bin/scen3/check_common_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_common_bus_interface]
command=./bin/scen3/common_bus_interface.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen3_check_internal_data_bus]
command=./bin/scen3/check_internal_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

