# supvisors/test/use_cases/scenario_3/etc/server/group_server.ini
[group:scen3_srv]
programs=item_manager,track_manager,common_bus_interface,check_common_data_bus,check_internal_data_bus

[program:item_manager]
command=./bin/item_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:track_manager]
command=./bin/track_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:check_common_data_bus]
command=./bin/check_common_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:common_bus_interface]
command=./bin/common_bus_interface.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:check_internal_data_bus]
command=./bin/check_internal_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

