# supvisors/test/etc/rocky51/service/program_flask.ini
[program:supvisorsflask]
command=supvisorsflask -u http://supv01:30000
autostart=true
autorestart=unexpected
startsecs=2
startretries=0
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s_out.log
stdout_logfile_maxbytes=10KB
stdout_logfile_backups=1
