# supvisors/tools/test/etc/common/database/program_movie_server.ini
[program:movie_server]
command = python ./scripts/process_app.py -n %(program_name)s
process_name = movie_server_%(process_num)02d
numprocs = 3
numprocs_start = 1
autostart = false
autorestart = false
startsecs = 1
startretries = 0
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

