# supvisors/tools/test/etc/common/player/program_movie_player.ini
[program:movie_player_01]
command = python ./scripts/process_app.py -n %(program_name)s
autostart = false
autorestart = false
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:movie_player_02]
command = python ./scripts/process_app.py -n %(program_name)s
autostart = false
autorestart = false
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:movie_player_03]
command = python ./scripts/process_app.py -n %(program_name)s
autostart = false
autorestart = false
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:movie_player_04]
command = python ./scripts/process_app.py -n %(program_name)s
autostart = false
autorestart = false
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

