# supvisors/test/use_cases/scenario_2/template_etc/server/group_server.ini
[group:scen2_services]
programs=config_manager,common_bus_interface,common_check_data_bus,internal_check_data_bus,internal_data_bus,data_processing,sensor_acquisition

[program:config_manager]
command=./bin/config_manager.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:common_check_data_bus]
command=./bin/common_check_data_bus.sh %(program_name)s
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:common_bus_interface]
command=./bin/common_bus_interface.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:internal_check_data_bus]
command=./bin/internal_check_data_bus.sh %(program_name)s
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:internal_data_bus]
command=./bin/internal_data_bus.sh %(program_name)s
autostart=false
startsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:data_processing]
command=./bin/data_processing.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:sensor_acquisition]
command=./bin/sensor_acquisition.sh %(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

