# supvisors/test/use_cases/scenario_2/template_etc/server/group_server.ini
[group:scen2_srv]
programs=scen2_config_manager,scen2_common_bus_interface,scen2_check_common_data_bus,scen2_check_internal_data_bus,scen2_internal_data_bus,scen2_data_processing,scen2_sensor_acquisition

[program:scen2_config_manager]
command=./bin/config_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_check_common_data_bus]
command=./bin/check_common_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_common_bus_interface]
command=./bin/common_bus_interface.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_check_internal_data_bus]
command=./bin/check_internal_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_internal_data_bus]
command=./bin/internal_data_bus.sh
autostart=false
startsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_data_processing]
command=./bin/data_processing.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_sensor_acquisition]
command=./bin/sensor_acquisition.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5
