# supvisors/test/use_cases/scenario_2/template_etc/server/programs_server.ini
[program:scen2_config_manager]
command=./bin/scen2/config_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_check_common_data_bus]
command=./bin/scen2/check_common_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_common_bus_interface]
command=./bin/scen2/common_bus_interface.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_check_internal_data_bus]
command=./bin/scen2/check_internal_data_bus.sh
autostart=false
startsecs=0
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_internal_data_bus]
command=./bin/scen2/internal_data_bus.sh
autostart=false
startsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_data_processing]
command=./bin/scen2/data_processing.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

[program:scen2_sensor_acquisition]
command=./bin/scen2/sensor_acquisition.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups=5

