# supvisors/test/use_cases/scenario_3/template_etc/common/group_common.ini
[group:scen3_common]
programs = internal_data_bus

[program:internal_data_bus]
command=./bin/internal_data_bus.sh
autostart=false
startsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_backups=5
