# supvisors/test/use_cases/scenario_2/etc/template/console/group_hci.ini
[group:scenario_2_hci]
programs=chart_view,sensor_control,sensor_view

[program:chart_view]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:sensor_control]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:sensor_view]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5
