# supvisors/test/use_cases/scenario_1/etc/cliche83/programs_cliche83.ini
[program:sensor_acquisition_2]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:sensor_processing_2]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5
