# supvisors/test/use_cases/scenario_1/etc/cliche83/wait_nfs_mount.ini
[program:scen1_wait_nfs_mount_3]
command=sleep 5
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5
