# supvisors/test/test_supervisor_177/etc/supervisord.conf
[inet_http_server]
port=:60000

[supervisord]
logfile=log/supervisord.log
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
umask=002

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:supvisors]
supervisor.rpcinterface_factory = supvisors.plugin:make_supvisors_rpcinterface

[supervisorctl]
serverurl=http://localhost:60000

[ctlplugin:supvisors]
supervisor.ctl_factory = supvisors.supvisorsctl:make_supvisors_controller_plugin


# Configuration to test Supervisor issue #177
[group:dummy1]
programs=autostarted,stopped

[group:dummy2]
programs=autostarted,control

[program:control]
command=zenity --info --text=%(program_name)s
autostart=true
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log

[program:autostarted]
command=zenity --info --text=%(group_name)s-%(program_name)s_%(process_num)d
process_name=%(program_name)s_%(process_num)d
numprocs=2
numprocs_start=10
autostart=true
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(group_name)s-%(program_name)s_%(process_num)d.log

[program:stopped]
command=zenity --text=%(group_name)s-%(program_name)s
process_name=%(program_name)s_%(process_num)d
numprocs=2
numprocs_start=20
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(process_num)d.log

[eventlistener:listener]
command=python bin/evt_listener.py
events=EVENT
buffer_size=50
autostart=true
stderr_logfile=log/%(program_name)s.log

[eventlistener:multiple_listener]
command=python bin/evt_listener.py
process_name=%(program_name)s_%(process_num)02d
numprocs=2
events=PROCESS_STATE
buffer_size=50
autostart=true
stderr_logfile=log/%(program_name)s_%(process_num)d.log
