# supvisors/test/use_cases/scenario_3/etc/console/cliche83/group_scen3_hci_01.ini
[group:scen3_hci_01]
programs = chart_view_01,item_control_01,check_internal_data_bus_01

[program:chart_view_01]
command = ./bin/chart_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:item_control_01]
command = ./bin/item_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_internal_data_bus_01]
command = ./bin/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

