# supvisors/tools/test/etc/cliche81/test/program_check_running_strategy.ini
[program:check_running_strategy]
command = python -m unittest scripts.check_running_strategy
autostart = false
autorestart = false
startsecs = 0
startretries = 0
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 10
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

