# supvisors/test/use_cases/scenario_2/etc/console/group_hci.ini
[group:scen2_hci_01]
programs = chart_view_01,sensor_control_01,sensor_view_01,hci_internal_check_data_bus_01,hci_internal_data_bus_01

[group:scen2_hci_02]
programs = chart_view_02,sensor_control_02,sensor_view_02,hci_internal_check_data_bus_02,hci_internal_data_bus_02

[group:scen2_hci_03]
programs = chart_view_03,sensor_control_03,sensor_view_03,hci_internal_check_data_bus_03,hci_internal_data_bus_03

[program:chart_view_01]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:chart_view_02]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:chart_view_03]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_control_01]
command = ./bin/sensor_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_control_02]
command = ./bin/sensor_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_control_03]
command = ./bin/sensor_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_view_01]
command = ./bin/sensor_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_view_02]
command = ./bin/sensor_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_view_03]
command = ./bin/sensor_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:hci_internal_check_data_bus_01]
command = ./bin/internal_check_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:hci_internal_check_data_bus_02]
command = ./bin/internal_check_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:hci_internal_check_data_bus_03]
command = ./bin/internal_check_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:hci_internal_data_bus_01]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:hci_internal_data_bus_02]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:hci_internal_data_bus_03]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

