# supvisors/test/use_cases/scenario_2/etc/console/group_scen2_hci_02.ini
[group:scen2_hci_02]
programs = scen2_chart_view_02,scen2_sensor_control_02,scen2_sensor_view_02,scen2_check_internal_data_bus_02,scen2_internal_data_bus_02

[program:scen2_chart_view_02]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_sensor_control_02]
command = ./bin/sensor_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_sensor_view_02]
command = ./bin/sensor_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_check_internal_data_bus_02]
command = ./bin/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_internal_data_bus_02]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

