# supvisors/test/use_cases/scenario_2/etc/console/programs_console.ini
[program:scen2_chart_view]
command = ./bin/scen2/chart_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_sensor_control]
command = ./bin/scen2/sensor_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_sensor_view]
command = ./bin/scen2/sensor_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_check_internal_data_bus]
command = ./bin/scen2/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups = 5

[program:scen2_internal_data_bus]
command = ./bin/scen2/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(group_name)s_%(program_name)s.log
stdout_logfile_backups = 5

