# supvisors/test/etc/common/import_database/program_mount_disk.ini
[program:mount_disk]
command=python ./scripts/process_app.py -n %(program_name)s
process_name=mount_disk_%(process_num)02d
numprocs=2
autostart=false
autorestart=false
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
