# supvisors/test/use_cases/scenario_2/etc/template/server/group_server.ini
[group:scenario_2_services]
programs=config_manager,load_data_model,data_bus,data_processing,sensor_acquisition

[program:config_manager]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:load_data_model]
command=zenity --info --text=%(program_name)s
autostart=false
startsecs=2
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:data_bus]
command=zenity --info --text=%(program_name)s
autostart=false
startsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:data_processing]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5

[program:sensor_acquisition]
command=zenity --info --text=%(program_name)s
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s.log
stdout_logfile_backups=5
