# supvisors/test/etc/common/player/program_test_reader.ini
[program:test_reader]
command=python ./scripts/process_app.py -n %(program_name)s -x 2
autostart=false
autorestart=false
startretries=0
exitcodes=5
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
