# supvisors/test/etc/import_database/program_copy_error.ini
[program:copy_error]
command=python ./scripts/process_app.py -n %(program_name)s -x 1
autostart=false
autorestart=false
startsecs=2
startretries=0
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
