# supvisors/test/use_cases/scenario_2/etc/server/group_server.ini
[group:scen2_services_01]
programs = config_manager_01,common_bus_interface_01,common_check_data_bus_01,internal_check_data_bus_01,internal_data_bus_01,data_processing_01,sensor_acquisition_01

[group:scen2_services_02]
programs = config_manager_02,common_bus_interface_02,common_check_data_bus_02,internal_check_data_bus_02,internal_data_bus_02,data_processing_02,sensor_acquisition_02

[program:config_manager_01]
command = ./bin/config_manager.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:config_manager_02]
command = ./bin/config_manager.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_bus_interface_01]
command = ./bin/common_bus_interface.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_bus_interface_02]
command = ./bin/common_bus_interface.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_check_data_bus_01]
command = ./bin/common_check_data_bus.sh %(program_name)s
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_check_data_bus_02]
command = ./bin/common_check_data_bus.sh %(program_name)s
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_check_data_bus_01]
command = ./bin/internal_check_data_bus.sh %(program_name)s
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_check_data_bus_02]
command = ./bin/internal_check_data_bus.sh %(program_name)s
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_data_bus_01]
command = ./bin/internal_data_bus.sh %(program_name)s
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_data_bus_02]
command = ./bin/internal_data_bus.sh %(program_name)s
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:data_processing_01]
command = ./bin/data_processing.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:data_processing_02]
command = ./bin/data_processing.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_acquisition_01]
command = ./bin/sensor_acquisition.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_acquisition_02]
command = ./bin/sensor_acquisition.sh %(program_name)s
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

