# supvisors/test/use_cases/scenario_2/etc/server/group_server.ini
[group:scen2_srv_01]
programs = config_manager_01,common_bus_interface_01,check_common_data_bus_01,check_internal_data_bus_01,internal_data_bus_01,data_processing_01,sensor_acquisition_01

[group:scen2_srv_02]
programs = config_manager_02,common_bus_interface_02,check_common_data_bus_02,check_internal_data_bus_02,internal_data_bus_02,data_processing_02,sensor_acquisition_02

[group:scen2_srv_03]
programs = config_manager_03,common_bus_interface_03,check_common_data_bus_03,check_internal_data_bus_03,internal_data_bus_03,data_processing_03,sensor_acquisition_03

[program:config_manager_01]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_bus_interface_01]
command = ./bin/common_bus_interface.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_common_data_bus_01]
command = ./bin/check_common_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_internal_data_bus_01]
command = ./bin/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_data_bus_01]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:data_processing_01]
command = ./bin/data_processing.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_acquisition_01]
command = ./bin/sensor_acquisition.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:config_manager_02]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_bus_interface_02]
command = ./bin/common_bus_interface.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_common_data_bus_02]
command = ./bin/check_common_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_internal_data_bus_02]
command = ./bin/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_data_bus_02]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:data_processing_02]
command = ./bin/data_processing.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_acquisition_02]
command = ./bin/sensor_acquisition.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:config_manager_03]
command = ./bin/config_manager.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:common_bus_interface_03]
command = ./bin/common_bus_interface.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_common_data_bus_03]
command = ./bin/check_common_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:check_internal_data_bus_03]
command = ./bin/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:internal_data_bus_03]
command = ./bin/internal_data_bus.sh
autostart = false
startsecs = 10
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:data_processing_03]
command = ./bin/data_processing.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:sensor_acquisition_03]
command = ./bin/sensor_acquisition.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

