# supvisors/test/etc/rocky51/test/program_check_stop_sequence.ini
[program:check_stop_sequence]
command=python -m unittest -c scripts.check_stop_sequence
autostart=false
autorestart=false
startsecs=0
startretries=0
exitcodes=0
stopsignal=INT
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
