# supvisors/test/use_cases/scenario_1/etc/rocky52/programs_rocky52.ini
[program:scen1_sensor_acquisition_1]
command=bin/scen1/sensor_acquisition.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:scen1_sensor_processing_1]
command=bin/scen1/sensor_processing.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

