# supvisors/test/etc/common/my_movies/program_manager.ini
[program:manager]
command=python ./scripts/process_app.py -n %(program_name)s
autostart=false
autorestart=unexpected
startsecs=0
startretries=0
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=./log/%(program_name)s_%(host_node_name)s_out.log
stdout_logfile_maxbytes=10KB
stdout_logfile_backups=1
stderr_logfile=./log/%(program_name)s_%(host_node_name)s_err.log
stderr_logfile_maxbytes=10KB
stderr_logfile_backups=1