# supvisors/test/etc/common/my_movies/program_converter.ini
[program:converter]
command=python ./scripts/process_app.py -n %(program_name)s -x 60
;command=bash -c 'python ./scripts/process_app.py -n %(program_name)s -x 10 && supvisorsctl -s http://localhost:60000 restart_process CONFIG %(group_name)s:%(program_name)s_%(process_num)02d'
process_name=converter_%(process_num)02d
numprocs=10
autostart=false
autorestart=false
startsecs=2
startretries=0
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=./log/%(program_name)s_%(process_num)02d_%(host_node_name)s_out.log
stdout_logfile_maxbytes=10KB
stdout_logfile_backups=1
stderr_logfile=./log/%(program_name)s_%(process_num)02d_%(host_node_name)s_err.log
stderr_logfile_maxbytes=10KB
stderr_logfile_backups=1
