# supvisors/test/use_cases/scenario_3/etc/console/cliche84/group_scen3_hci_02.ini
[group:scen3_hci_02]
programs = scen3_chart_view_02,scen3_item_control_02,scen3_check_internal_data_bus_02

[program:scen3_chart_view_02]
command = ./bin/chart_view.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:scen3_item_control_02]
command = ./bin/item_control.sh
autostart = false
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

[program:scen3_check_internal_data_bus_02]
command = ./bin/check_internal_data_bus.sh
autostart = false
startsecs = 0
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s.log
stdout_logfile_backups = 5

