# supvisors/test/use_cases/scenario_2/etc/common/programs_services.ini
[program:common_data_bus]
command = zenity --info --text=%(program_name)s
startsecs = 40
killasgroup = true
stopasgroup = true
redirect_stderr = true
stdout_logfile = log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_backups = 5
