# supvisors/tools/test/etc/common/web_movies/program_web_browser.ini
[program:web_browser_01]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_02]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_03]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_04]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_05]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_06]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_07]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_08]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_09]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_10]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_11]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_12]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_13]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_14]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_15]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

[program:web_browser_16]
command = firefox --new-window http://localhost:60000
autostart = false
autorestart = false
startsecs = 5
startretries = 0
redirect_stderr = true
stdout_logfile = ./log/%(program_name)s_%(host_node_name)s.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1

