# supvisors/test/use_cases/scenario_1/etc/cliche82/group-scenario_1.ini
[group:scenario_1]
programs=wait_nfs_mount_2,sensor_acquisition_1,sensor_processing_1

[program:wait_nfs_mount_2]
command=./bin/wait_nfs_mount.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:sensor_acquisition_1]
command=./bin/sensor_acquisition.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:sensor_processing_1]
command=./bin/sensor_processing.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

