# supvisors/test/use_cases/scenario_1/etc/cliche81/group-scenario_1.ini
[group:scenario_1]
programs=hci,config_manager,data_processing,external_interface,data_recorder,wait_nfs_mount_1

[program:hci]
command=./bin/hci.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:config_manager]
command=./bin/config_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:data_processing]
command=./bin/data_processing.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:external_interface]
command=./bin/external_interface.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:data_recorder]
command=./bin/data_recorder.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:wait_nfs_mount_1]
command=./bin/wait_nfs_mount.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

