# supvisors/test/use_cases/scenario_1/etc/cliche81/programs_cliche81.ini
[program:hci]
command=zenity --info --text=%(program_name)s
autostart=false
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:config_manager]
command=zenity --info --text=%(program_name)s
autostart=false
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:data_processing]
command=zenity --info --text=%(program_name)s
autostart=false
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:external_interface]
command=zenity --info --text=%(program_name)s
autostart=false
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:data_recorder]
command=zenity --info --text=%(program_name)s
autostart=false
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5
