# supvisors/test/use_cases/scenario_1/etc/rocky51/programs_rocky51.ini
[program:scen1_hci]
command=bin/scen1/hci.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:scen1_config_manager]
command=bin/scen1/config_manager.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:scen1_data_processing]
command=bin/scen1/data_processing.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:scen1_external_interface]
command=bin/scen1/external_interface.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5

[program:scen1_data_recorder]
command=bin/scen1/data_recorder.sh
autostart=false
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=log/%(program_name)s_%(ENV_CUR_DATE)s_%(ENV_CUR_TIME)s.log
stdout_logfile_backups=5
